# services/docker_compose_manager/dockcoman.timer
[Unit]
Description=docker-compose manager timer

[Timer]
OnUnitActiveSec=20s
OnBootSec=20s

[Install]
WantedBy=timers.target